# gigi.service
[Unit]
Description=Gigi Discord Bot
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=YOUR_APP_DIR
EnvironmentFile=YOUR_APP_DIR/.env
ExecStart=/usr/bin/node YOUR_APP_DIR/dist/bot.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
